# cleanup-repo-tools/cleanup-repo-tool.timer
[Unit]
Description=Run Local Repository Cleanup Daily

[Timer]
#OnCalendar=daily
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
